# appliance/systemd/haai-api.service
[Unit]
Description=Home Assistant AI appliance API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=haai
Group=haai
WorkingDirectory=/opt/haai
EnvironmentFile=-/etc/haai/haai.env
Environment=HAAI_HOST=0.0.0.0
Environment=HAAI_DATA_DIR=/var/lib/haai
ExecStart=/usr/bin/node /opt/haai/dist/server/index.js
Restart=always
RestartSec=5
# Allow the non-root haai user to bind standard web ports such as 80 and 443.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
# The API can only start the root updater through the narrow sudoers rule in
# /etc/sudoers.d/haai-updater. NoNewPrivileges would block that sudo handoff.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/haai /var/log/haai

[Install]
WantedBy=multi-user.target
